# systemd/nwingest.service
# NightWatcher FITS ingest watcher.
#
# Installed by the nightwatcher-ingest package to /lib/systemd/system/.
# Configure /etc/nwingest/nwingest.yaml (and set NWDB_PASSWORD in
# /etc/nwingest/nwingest.env if you enable the SQM stamp or the web-UI
# extension), then:  sudo systemctl enable --now nwingest
# Watch it:          journalctl -fu nwingest

[Unit]
Description=NightWatcher FITS ingest watcher
After=network-online.target remote-fs.target
Wants=network-online.target

[Service]
Type=simple
# The service account must read incoming/ and write the archive tree (grant via
# group membership or an ACL). To run as a different user, add a drop-in:
#   systemctl edit nwingest   ->   [Service]\nUser=youruser\nGroup=yourgroup
User=nwingest
Group=nwingest
ExecStart=/usr/bin/nwingest --config /etc/nwingest/nwingest.yaml watch
Restart=on-failure
RestartSec=10
# Secrets (e.g. NWDB_PASSWORD=...) live here, not in the unit or the YAML.
EnvironmentFile=-/etc/nwingest/nwingest.env

[Install]
WantedBy=multi-user.target
